# debian/postrm
#!/bin/bash
set -e

case "$1" in
    purge)
        # Remove user and group
        if getent passwd rustbox >/dev/null; then
            deluser rustbox
        fi
        if getent group rustbox >/dev/null; then
            delgroup rustbox
        fi
        
        # Remove directories
        rm -rf /var/lib/rustbox
        rm -rf /var/log/rustbox
        rm -rf /tmp/rustbox-locks
        
        # Remove systemd units
        systemctl stop rustbox.service 2>/dev/null || true
        systemctl disable rustbox.service 2>/dev/null || true
        systemctl daemon-reload
        ;;
    remove|upgrade|failed-upgrade|abort-install|abort-upgrade|disappear)
        # Stop service but keep data
        systemctl stop rustbox.service 2>/dev/null || true
        ;;
esac

exit 0